# systemd-access-control/access-control.service
# /etc/systemd/system/access-control.service
[Unit]
Description=Access control via MQTT and MySQL
After=network-online.target mosquitto.service mysql.service
Requires=mosquitto.service mysql.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/systemd-access-control
ExecStart=/root/systemd-access-control/venv/bin/python3 -u /root/systemd-access-control/access-control.py
Environment="PATH=/root/systemd-access-control/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin"
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target